# Image_Repository/supervisord/ini/laravel-queue.ini
[program:laravel-worker]

process_name=%(program_name)s_%(process_num)02d

command=/usr/bin/php /var/www/html/artisan queue:listen --sleep=3 --tries=3 --timeout=0

autostart=true
autorestart=true

numprocs=4

redirect_stderr=true